# deploy/systemd/stockwatch-all-in-one.service
[Unit]
Description=StockWatch All-in-One Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kac0
Group=kac0
WorkingDirectory=/opt/stockwatch
ExecStart=/opt/stockwatch/stockwatchctl all-in-one
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
