# deploy/systemd/stockwatch-ops.local.service
[Unit]
Description=StockWatch Ops Service (Local Host)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kac0
Group=kali
WorkingDirectory=/home/kac0/project/stockwatch
ExecStart=/home/kac0/project/stockwatch/stockwatchctl ops
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
